# scripts/unifi-bridge/unifi-bridge.service
[Unit]
Description=Glance UniFi Protect bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=unifi-bridge
WorkingDirectory=/opt/glance-unifi-bridge
EnvironmentFile=/opt/glance-unifi-bridge/.env
ExecStart=/usr/bin/node bridge.mjs
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
